# .config/autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Name=Network Manager Applet
Exec=nm-applet
